# infrastructure-provisioning/src/deeplearning/templates/jupyter-notebook.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

exec sudo -i -u OS_USR /bin/bash -c "/usr/local/bin/jupyter notebook --config CONF_PATH"