# infrastructure-provisioning/src/deeplearning/templates/tensorboard.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

exec sudo -i -u OS_USR python2.7 -m tensorflow.tensorboard --logdir=/var/log/tensorboard --port 6006