# infrastructure-provisioning/src/general/templates/os/inactive.service
[Unit]
Description=inactivity check

[Service]
User=OS_USR
Type=oneshot
ExecStart=/bin/bash /opt/inactivity/inactive.sh

