# infrastructure-provisioning/src/general/templates/os/inactive.timer
[Unit]
Description=inactive service

[Timer]
OnUnitActiveSec=3s
OnBootSec=120s
AccuracySec=5s

[Install]
WantedBy=timers.target